# Logs/starter.conf
[program:starter]
command=/bin/bash -c "/tmp/starter.sh"
autorestart=false
startsecs=0